# debian/vyatta-debian-system-config.postinst
#!/bin/sh
# postinst script for vyatta-debian-system-config
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

sysconfdir=/opt/vyatta/etc
pkg=vyatta-debian-system-config

case "$1" in
    configure)
	# init of daemon is controlled by Vyatta configuration
	for init in ntp snmpd openhpid logd isc-dhcp-relay \
            isc-dhcp-server ipvsadm ddclient radvd hostapd
	do
	    update-rc.d -f ${init} remove >/dev/null
	done

	# disable dnsmasq start-up. Controlled by Vyatta configuration.
	update-rc.d -f dnsmasq remove >/dev/null
	deb-systemd-helper disable dnsmasq.service

	# Udev package asks for user 'tss' early in boot process. Want to avoid
	# going out to remote services to look for this local user
	if ! grep -q '^tss' /etc/passwd; then
	    adduser --system --group --shell /usr/sbin/nologin \
		--home /var/lib/tpm tss
	fi

	touch /etc/sudoers
	cp -p /etc/sudoers /etc/sudoers.bak

	# Turn off Debian default for %sudo
	sed -i -e '/^%sudo/d' /etc/sudoers || true


	touch /etc/environment

	# Change default shell for new accounts
	sed -i -e ':^DSHELL:s:/bin/bash:/bin/vbash:' /etc/adduser.conf

	# Do not allow users to change full name field (controlled by Vyatta
	# config)
	sed -i -e 's/^CHFN_RESTRICT/#&/' /etc/login.defs
	# Set default umask
	sed -i -e 's/^UMASK.*/UMASK 027/' /etc/login.defs

	# Disable reboot on ctrl-alt-delete
	if [ -e /etc/inittab ]; then
		# Debian 7 sysvinit specific
		sed -i -e '/^ca/s/\/sbin\/shutdown .*$/\/bin\/echo "Ctrl\-Alt\-Del disabled."/' /etc/inittab

		# former lb hooks / 15-disable-TO.chroot
		sed -i '/^T0\:23/s/^/#/g' /etc/inittab
	else
		ln -sf /dev/null /etc/systemd/system/ctrl-alt-del.target
	fi

	# remove unnecessary ddclient script in /etc/ppp/ip-up.d/
	# this logs unnecessary messages trying to start ddclient
	rm -f /etc/ppp/ip-up.d/ddclient

	# remove the UNCONFIGURED in /etc/fstab
	sed -i -e '/^# UNCONFIGURED FSTAB FOR BASE SYSTEM/d' /etc/fstab

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
